# provisioning/roles/common/templates/etc/profile.d/terrarium.sh
export DATABASE_URL="postgresql://{{dbuser}}:{{dbpassword}}@localhost/{{dbname}}"
export TERRARIUM_SECRET_KEY={{TERRARIUM_SECRET_KEY}}
